# grading-app-celery-beat.service
[Unit]
Description=Document Grading App - Celery Beat Scheduler
After=network.target redis.service valkey.service
Wants=redis.service valkey.service

[Service]
Type=simple
User=grader
Group=grader
WorkingDirectory=/opt/grading-app
Environment=PATH=/opt/grading-app/venv/bin
ExecStart=/opt/grading-app/venv/bin/celery -A tasks beat --loglevel=info
Restart=always
RestartSec=10
KillSignal=SIGTERM

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/grading-app

[Install]
WantedBy=multi-user.target